# source/etc/init/app.conf
description "Gunicorn application server running app"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid pi
setgid www-data

